# templet/templet/deploy.bak/supervisor_commands/pollqueue.ini
[program:pollqueue]
command=/home/diors/env/bin/python /home/diors/work/templet/manage.py poll_queue
autostart=true
autorestart=true
startretries=3
exitcodes=0,1,2
stopsignal=KILL
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor_pollqueue_access.log
stderr_logfile=/var/log/supervisor_pollqueue_error.log
user=diors