# templet/templet/deploy.bak/supervisor_commands/rlistlog.ini
[program:rlistlog]
command=/home/diors/env/bin/python /home/diors/work/templet/manage.py rlistlog --key=django:logit:templet --filename=/tmp/templet.logit.log
autostart=true
autorestart=true
startretries=3
exitcodes=0,1,2
stopsignal=KILL
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor_rlistlog_access.log
stderr_logfile=/var/log/supervisor_rlistlog_error.log
user=diors