# templet/templet/deploy.bak/templet_supervisor2.ini
[program:templet2]
command=/home/diors/env/bin/uwsgi --ini /home/diors/work/templet/templet/deploy/templet2.ini
autostart=true
autorestart=true
startretries=3
exitcodes=0,1,2
stopsignal=KILL
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor_templet2_access.log
stderr_logfile=/var/log/supervisor_templet2_error.log
user=diors